# supervisord.conf
[supervisord]
nodaemon = true

[program:servicestack]
command = /usr/local/bin/mono /build/monotest.exe
autorestart = true
logfile = /var/log/supervisor/servicestack.log
log_stderr = true
